# development/conf_supervisor_simple.conf
[unix_http_server]
file=%(here)s/supervisor/supervisor_simple.sock

[supervisord]
logfile=%(here)s/supervisor/supervisord_simple.log
loglevel=debug
pidfile=%(here)s/supervisor/supervisord_simple.pid
childlogdir=%(here)s/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://%(here)s/supervisor/supervisor_simple.sock

[program:distributed-worker-8001]
directory=%(here)s
command=poetry run python -m distributed.simple.worker --path=%(here)s/supervisor/simple_worker_8001.sock --port 8001
priority=900
stopasgroup=true
killasgroup=true

[program:distributed-worker-8002]
directory=%(here)s
command=poetry run python -m distributed.simple.worker --path=%(here)s/supervisor/simple_worker_8002.sock --port 8002
priority=900
stopasgroup=true
killasgroup=true

[program:distributed-worker-8003]
directory=%(here)s
command=poetry run python -m distributed.simple.worker --path=%(here)s/supervisor/simple_worker_8003.sock --port 8003
priority=900
stopasgroup=true
killasgroup=true

[program:distributed-worker-8004]
directory=%(here)s
command=poetry run python -m distributed.simple.worker --path=%(here)s/supervisor/simple_worker_8004.sock --port 8004
priority=900
stopasgroup=true
killasgroup=true

[group:distributed-workers]
programs=distributed-worker-8001,distributed-worker-8002,distributed-worker-8003,distributed-worker-8004
priority=900

[program:distributed-master]
directory=%(here)s
process_name=%(program_name)s
command=poetry run python -m distributed.simple.master
priority=999
autostart=false
autorestart=false
